# agent-python/ransomshield-agent.service
[Unit]
Description=RansomShield Host Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/opt/ransomshield-agent
ExecStart=/opt/ransomshield-agent/venv/bin/python /opt/ransomshield-agent/ransomshield_host_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
