# agent-python/ransomshield-agent.service
[Unit]
Description=RansomShield Host Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ransomshield
EnvironmentFile=/opt/ransomshield/.env
ExecStart=/opt/ransomshield/venv/bin/python /opt/ransomshield/ransomshield_host_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
